# patches/etc/init.d/set-ocs-hostname
#!/sbin/runscript
# Copyright (c) 2015 Online Labs Cloud Team <cloud-team@labs.online.net>
# Released under the MIT license.

description="Set OCS hostname"

depend() {
    need root
}

start() {
    ebegin "setting OCS hostname"
    echo "hostname=\"$(cat /etc/hostname)\"" > /etc/conf.d/hostname
    eend $?
}
