# patches/etc/init.d/sync-kernel-extra
#!/sbin/runscript
# Copyright (c) 2015 Online Labs Cloud Team <cloud-team@labs.online.net>
# Released under the MIT license.

description="Downloads kernel modules/headers from OCS' mirror"

depend() {
    need root
    before kmod-static-nodes
}

start() {
    ebegin "Syncing kernel modules/headers"
    /usr/local/sbin/oc-sync-kernel-headers && /usr/local/sbin/oc-sync-kernel-modules
    eend $?
}
